# scripts/cloudflare-tunnel.service.template
[Unit]
Description=Cloudflare Tunnel for note-mcp-server
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
ExecStart=/usr/local/bin/cloudflared tunnel run --loglevel debug note-mcp
Restart=always
RestartSec=10
WorkingDirectory=/home/YOUR_USERNAME/noteMCP
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
