# scripts/cloudflare-tunnel.service
[Unit]
Description=Cloudflare Tunnel for note-mcp-server
After=network.target

[Service]
Type=simple
User=***USERNAME_REMOVED***
ExecStart=/usr/local/bin/cloudflared tunnel run --loglevel debug note-mcp
Restart=always
RestartSec=10
WorkingDirectory=/Users/***USERNAME_REMOVED***/noteMCP
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
